# iac-workspace/ansible/roles/web/templates/ssl-cert-check.service.j2
[Unit]
Description=SSL Certificate Expiry Check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ssl-cert-check.sh
